# systemd/diamond.service
[Unit]
Description=Diamond Learning Server
After=network.target

[Service]
Type=simple
User=diamond
WorkingDirectory=/opt/diamond
ExecStart=/opt/diamond/diamond
Restart=on-failure
RestartSec=5

Environment=DIAMOND_LLM_URL=https://api.deepseek.com
Environment=DIAMOND_LLM_KEY=your-deepseek-api-key-here
Environment=DIAMOND_MODEL=deepseek-chat
Environment=DIAMOND_PORT=7331
Environment=DIAMOND_DATA_DIR=/var/lib/diamond

[Install]
WantedBy=multi-user.target
